# crontab
# Запуск парсера каждые 4 часа
0 */4 * * * /app/run_parser.sh >> /var/log/cron.log 2>&1
# Пустая строка в конце обязательна для cron